# crontab
@reboot sh -c 'cd ~/letsmoe.com && pm2 start server.mjs'